# smartthings_oven_bridge/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv sh
set -eu

cd /
exec python3 -m app.main
